# rootfs/overlay/etc/sv/platform-mgrd/run
#!/bin/sh
exec /usr/sbin/platform-mgrd 2>&1
