# rootfs/overlay/etc/systemd/system/platform-mgrd.service
[Unit]
Description=NOS Platform Manager (fan/temp/LEDs)
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/sbin/platform-mgrd
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
